# scraper/crontab

0 4 * * * cd cd /home/jlogan/job_tracker/scraper && /home/jlogan/.venv/bin/python build_daily_snapshot.py